# containerized/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
user=nginx
command=/usr/sbin/php-fpm7 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
